# no.prikkprikkprikk.obsidian.backup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>no.prikkprikkprikk.obsidian.backup</string>
    <key>ProgramArguments</key>
    <array>
      <string>/bin/sh</string>
      <string>-c</string>
      <string>/usr/bin/shortcuts run "Obsidian Git Backup"</string>
    </array>
    <key>StartInterval</key>
    <integer>3600</integer>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/prikkprikkprikk/Library/Logs/obsidian-backup.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/prikkprikkprikk/Library/Logs/obsidian-backup-error.log</string>
</dict>
</plist>
